# deployment/mybuddy-backend.service
[Unit]
Description=MyBuddy Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=mybuddy
Group=mybuddy
WorkingDirectory=/opt/mybuddy/backend
ExecStart=/usr/bin/java -jar -Dspring.profiles.active=prod -Dserver.address=127.0.0.1 -Dserver.port=8080 mybuddy-backend.jar
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mybuddy-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mybuddy/backend/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment variables
Environment="JAVA_OPTS=-Xms512m -Xmx1024m -XX:+UseG1GC -XX:+UseStringDeduplication"

[Install]
WantedBy=multi-user.target
